# scripts/systemd/paperclip-uimodifi-proxy.service
[Unit]
# Public :3301 must forward only to 127.0.0.1:3302 (preview app). A stray socat on :3301
# pointing elsewhere (e.g. :3003) steals the port; systemd then cannot bind and curl shows "Empty reply".
Description=Paperclip PaperclipUimodifi preview proxy (public :3301 to 127.0.0.1:3302)
After=network.target paperclip-uimodifi-preview.service

[Service]
User=ubuntu
ExecStart=/usr/bin/socat TCP4-LISTEN:3301,fork,reuseaddr TCP4:127.0.0.1:3302
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=paperclip-uimodifi-proxy

[Install]
WantedBy=multi-user.target
